# main/spamass-milter/spamass-milter.initd
#!/sbin/openrc-run

command=/usr/sbin/spamass-milter
command_args=${command_args:--p /run/spamass-milter.sock}
pidfile=${pidfile:-/run/spamass-milter.pid}

depend() {
	need net
	before mta
}

start_pre() {
	[ -n "$PIDFILE" ] && pidfile=$PIDFILE
	[ -n "$EXTRA_OPTS" ] && command_args=$EXTRA_OPTS

	command_args="-f -P $pidfile $command_args"
}
